# deploy/web-server-deploy/deploy/karaoke-docker-prune.timer
[Unit]
Description=Karaoke: ежедневная очистка неиспользуемых Docker-образов

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
